# docker/supervisor/conf.d/event_sample.conf
[program:event]
command=php /src/run.php
autostart=false
autorestart=false
startretries=0
user=root
numprocs=1
redirect_stderr=true
stdout_logfile=/etc/supervisor/log/event_sample.log
